# vikunja-backup.timer
[Unit]
Description=Daily Vikunja backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
